# root/s6/debian-root/etc/services.d/pihole-FTL/finish
#!/usr/bin/with-contenv bash

PID=$(pgrep pihole-FTL)

if [ -n "$PID" ]; then
  s6-echo "Stopping pihole-FTL"
  kill -9 $PID
else
  s6-echo "Stopping pihole-FTL (not running)"
fi
